# bin/jupyter.service
[Unit]
Description=Jupyter Notebook Server

[Service]
PIDFILE=/run/jupyter.pid
User=aimmodev
Environment="PYTHONPATH=/home/aimmodev/jupyter-workspace/venv/lib/python3.10/site-packages"
Environmant="PATH=/home/aimmodev/jupyter-workspace/venv/bin:/bin:/usr/bin:/usr/local/bin:/sbin:/usr/sbin:/usr/local/sbin"
WorkingDirectory=/home/aimmodev/jupyter-workspace
ExecStart=/home/aimmodev/jupyter-workspace/venv/bin/jupyter-lab

[Install]
WantedBy=multi-user.target
